# zcu111d/conf/io.piradio.zcu111d.service
[Unit]
Description=zcu111d

[Service]
ExecStart=/usr/bin/zcu111d

[Install]
WantedBy=multi-user.target