# piradio-fpgad/conf/piradio-fpgad.service
[Unit]
Description=Pi Radio FPGA Management Daemon

[Service]
ExecStart=/usr/bin/piradio-fpgad

[Install]
WantedBy=multi-user.target