# rfdcd/conf/io.piradio.rfdcd.service
[Unit]
Description=rfdcd

[Service]
ExecStart=/usr/bin/piraido-rfdcd

[Install]
Alias=piradio-rfdcd.service
WantedBy=multi-user.target