# piradio-sampled/conf/piradio-sampled.service
[Unit]
Description=Pi Radio Non Real Time sample daemon

[Service]
ExecStart=/usr/bin/piradio-sampled

#[Install]
#WantedBy=multi-user.target