# sampled/deb/postinst
#!/bin/bash

set -e

deb-systemd-helper enable io.piradio.sampled.service
